# .dotfiles/systemd/btrfs_snap.service
[Unit]
Description=Snap Settimanale

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /home/bcclsn/.script/btrfs_snap.sh 
