# .script/systemd/backup.timer
[Unit]
Description=Backup Giornaliero

[Timer]
OnCalendar=*-*-* 12:00:00
Unit=backup.service
Persistent=true

[Install]
WantedBy=timers.target
