# .script/systemd/backup.service
[Unit]
Description=Backup Settimanale

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /home/bcclsn/.script/backup_b2.sh 
